# broadlink-ac-mqtt/rootfs/etc/services.d/broadlink-ac-mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

python3 -c 'import sys, yaml, json; print(yaml.dump(json.loads(open("/data/options.json").read())))' > /usr/share/broadlink_ac_mqtt-1.2.4/config.yml

## Check if debug mode is enabled
DEBUG_FLAG=""
if bashio::config.true 'service.debug'; then
    export BROADLINK_DEBUG=true
    DEBUG_FLAG="-d"
    bashio::log.info "Debug mode is enabled — verbose logging is active"
fi

## Forward log_level / log_level_console from add-on config to monitor.py env vars.
## Ignored by monitor.py if --debug is also passed (debug wins).
if bashio::config.has_value 'service.log_level'; then
    export LOG_LEVEL="$(bashio::config 'service.log_level')"
    bashio::log.info "LOG_LEVEL=${LOG_LEVEL}"
fi
if bashio::config.has_value 'service.log_level_console'; then
    export LOG_LEVEL_CONSOLE="$(bashio::config 'service.log_level_console')"
    bashio::log.info "LOG_LEVEL_CONSOLE=${LOG_LEVEL_CONSOLE}"
fi

## Run your program
python3 /usr/share/broadlink_ac_mqtt-1.2.4/monitor.py ${DEBUG_FLAG}
